# debian autostart/avsystem.desktop
[Desktop Entry]
Type=Application
Name=AVSystem
Exec=bash -c 'cd /home/prg/Extron/ServerLib/source && ./run_in_terminal.sh'
Terminal=true;
